# system/etc/udev/rules.d/85-my_usb_device_rule.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/home/me/bin/autonumlock"
ACTION=="remove", SUBSYSTEM=="usb", RUN+="/home/me/bin/autonumlock"
